# crontab
*/5 * * * *	munin	/etc/munin/huawei-hg612-munin/getstats.sh
